# static-modules/oddb_org/templates/currency.erb
#!/sbin/runscript

CURRENCY_BINARY=<%= scope.lookupvar('currency_installed') %>
CURRENCY_PIDFILE=${CURRENCY_PIDFILE:-/var/run/${SVCNAME}.pid}

checkconfig() {
  if [ ! -x $CURRENCY_BINARY ] ; then
    return 1
  fi
  eselect ruby set ruby19
  return 0
}

start() {
  checkconfig || return 1

  ebegin "Starting ${SVCNAME}"
  start-stop-daemon --start --exec "${CURRENCY_BINARY}" \
      --pidfile "${CURRENCY_PIDFILE}" --make-pidfile  --background
  eend $?
}

stop() {
  if [ "${RC_CMD}" = "restart" ] ; then
    checkconfig || return 1
  fi

  ebegin "Stopping ${SVCNAME}"
  start-stop-daemon --stop --exec "${CURRENCY_BINARY}" \
      --pidfile "${CURRENCY_PIDFILE}" --quiet
  eend $?
}

reload() {
  checkconfig || return 1
  ebegin "Reloading ${SVCNAME}"
  start-stop-daemon --signal HUP \
      --exec "${CURRENCY_BINARY}" --pidfile "${CURRENCY_PIDFILE}"
  eend $?
}
