# static-modules/oddb_org/templates/service.erb
#!/sbin/runscript
# KEYWORD: shutdown
<% if scope.lookupvar('service_user') %>
SERVICE_USER="--user <%= scope.lookupvar('service_user') %>"
<% end %>
SERVICE_BINARY="<%= scope.lookupvar('service_location') %>"
SERVICE_DEPEND="<%= scope.lookupvar('service_depend') %>"
SERVICE_PIDFILE=${SERVICE_PIDFILE:-/var/run/${SVCNAME}.pid}

checkconfig() {
  if [ ! -x $SERVICE_BINARY ] ; then
    return 1
  fi
  return 0
}

depend() {
  need net
}

start() {
  checkconfig || return 1

  ebegin "Starting ${SVCNAME}"
  start-stop-daemon --start --exec "${SERVICE_BINARY}" ${SERVICE_USER} \
      --pidfile "${SERVICE_PIDFILE}" --make-pidfile  --background
  eend $?
}

stop() {
  if [ "${RC_CMD}" = "restart" ] ; then
    checkconfig || return 1
  fi

  ebegin "Stopping ${SVCNAME}"
  start-stop-daemon --stop --exec "${SERVICE_BINARY}" ${SERVICE_USER} \
      --pidfile "${SERVICE_PIDFILE}" --quiet
  eend $?
}

reload() {
  checkconfig || return 1
  ebegin "Reloading ${SVCNAME}"
  start-stop-daemon --signal HUP \
      --exec "${SERVICE_BINARY}" ${SERVICE_USER} --pidfile "${SERVICE_PIDFILE}"
  eend $?
}
